# test/e2e/infra/vagrant/playbook/roles/control-plane/templates/socat-dns-proxy-tcp.service.j2
[Unit]
Description=Socat DNS proxy TCP

[Service]
ExecStart=/usr/bin/socat TCP6-LISTEN:5353,fork TCP4-CONNECT:127.0.0.53:53

[Install]
WantedBy=multi-user.target
